# templates/default/kibana.service.erb
[Unit]
Description=<%= @service_name %>
After=network.target

[Service]
ExecStart=<%= @bindir %>/kibana  <%= @args %>
User=<%= @user %>
Restart=always
